# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=bash /app/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log